# cfg/core/systemd/xorg@.target
[Unit]
Description=Ready xorg target on display %i
Requires=xorg@%i.service
