# cfg/core/udev/70-wifi-powersave.rules
# To enable powersaving on all wireless interfaces:
ACTION=="add", SUBSYSTEM=="net", KERNEL=="wlan*", RUN+="/usr/bin/iw dev %k set power_save on"
