# cfg/core/udev/hdmi_monitor.rules
SUBSYSTEM=="drm", ACTION=="change", RUN+="/home/amerlyq/.bin/setX d"
