# cfg/xorg/xkb/unit/kbdd.service
[Unit]
Description=Xorg Keyboard PerWindow - %P
Requires=wm.target
After=wm.target
PartOf=main.target
Before=main.target

[Service]
# Type=dbus
# BusName=ru.gentoo.KbddService
ExecStart=/usr/bin/kbdd --nodaemon
# User=%i
# Group=users
# Restart=on-failure
Restart=always
RestartSec=1

[Install]
# Alias=dbus-ru.gentoo.KbddService.service
WantedBy=main.target
