# cfg/xorg/xkb/unit/xkb.service
[Unit]
Description=Xorg Keyboard Layout
PartOf=xorg.target
After=xorg.service
ConditionFileNotEmpty=%h/.cache/airy/keylayout.xkm
# OR ExecStart=-/... to suppress warnings

[Service]
# CHECK:BUG: oneshot is not restarted on 'start main.target'?
Type=oneshot
RemainAfterExit=yes
# BUG: ExecStart=/usr/bin/xkbcomp %h/.cache/airy/keylayout.xkm ${DISPLAY}
ExecStart=/bin/sh -c r.xkb-load
ExecStop=/usr/bin/setxkbmap -option ""

[Install]
# MAYBE: Impossible to be part of {xorg,wm}.target?
# CHECK: r.init-graph
WantedBy=xorg.target
