# xorg/unit/wm.target
[Unit]
Description=WM - (Target)
Conflicts=logoff.target
# PartOf=main.target
AllowIsolate=true

[Install]
RequiredBy=main.target
