# cfg/core/xorg/unit/xorg.target
[Unit]
Description=Xorg - (Target)
Conflicts=logoff.target
AllowIsolate=true
