# cfg/shell/profile.d/default
# vim: ft=sh

### General ###
export CURR_PROF="default"

# Set your language environment
export LANG=en_US.UTF-8
export LC_NUMERIC=en_US.UTF-8
# LANGUAGE=en_US

export EDITOR=vim
export VISUAL=vim

for nm in host dev
do source ~/.config/airy/$nm || return
done
