# cfg/shell/profile.d/ssh
# vim: ft=sh
source ~/.config/airy/default || return

### General ###
CURR_PROF="ssh"

### HW ###
CURR_PKG_LIST="shell dev net cli_utils ranger vim zsh"
