# cfg/core/init/unit/demons.target
[Unit]
Description=Demon - (Target)

[Install]
WantedBy=default.target
