# cfg/core/init/unit/demons.target
[Unit]
Description=Demon - (Target)
# After=default.target

[Install]
WantedBy=default.target
