# cfg/core/init/unit/main.target
[Unit]
Description=Xorg Demons
Requires=dbus.socket
Requires=xorg.service
Wants=wm.target
AllowIsolate=true
