# cfg/core/init/unit/main.target
[Unit]
Description=Xorg Desktop Environment
After=xorg.target
Wants=wm.target
After=wm.target
AllowIsolate=true
