# cfg/core/init/unit/wm.target
[Unit]
Description=Xorg WM %u@%H
After=xorg.target
PartOf=main.target
AllowIsolate=true

[Install]
WantedBy=main.target
