# cfg/core/init/unit/wm.target
[Unit]
Description=Xorg WM %u@%H
Requires=xorg.target
After=xorg.target
PartOf=main.target
Before=main.target
AllowIsolate=true
