# cfg/core/init/unit/wm.target
[Unit]
Description=Window Manager
Requires=xorg.socket
Requires=dbus.socket
Before=main.target
PartOf=main.target
AllowIsolate=true
