# cfg/core/init/unit/wm.target
[Unit]
Description=Window Manager for %u on %H
Requires=xorg.target
After=xorg.target
PartOf=main.target
Before=main.target
AllowIsolate=true
