# cfg/core/init/unit/xorg.target
[Unit]
Description=Xorg Target
Requires=dbus.socket
Requires=xorg.service
PartOf=main.target
Before=main.target
AllowIsolate=true
