# cfg/core/systemd/#/try.service
[Unit]
Description=Try notify

[Service]
Type=notify
NotifyAccess=all
ExecStart=/bin/sh -c 'bash -c "sleep 2; systemd-notify --ready"'
